# linux/systemd/denon-remap.service
[Unit]
Description=Load DENON 5.1 remap sink and HDMI profile
After=pipewire-pulse.service
Wants=pipewire-pulse.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -lc '\
  pactl set-card-profile alsa_card.pci-0000_00_1f.3 output:hdmi-surround && \
  pactl load-module module-remap-sink \
    sink_name=denon51_remap \
    master=alsa_output.pci-0000_00_1f.3.hdmi-surround \
    channels=6 \
    master_channel_map=front-left,front-right,rear-left,rear-right,front-center,lfe \
    channel_map=front-left,front-right,front-center,lfe,rear-left,rear-right \
    remix=yes \
'

[Install]
WantedBy=default.target
