# linux/systemd/htpc-update.timer
[Unit]
Description=Daily htpc-remote update at 2am

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
