# linux/systemd/ydotoold.service
[Unit]
Description=ydotool daemon for Wayland automation
After=graphical-session.target

[Service]
ExecStart=/usr/local/bin/ydotoold
Restart=always

[Install]
WantedBy=default.target
